# submit/prj1-sol/test/crontab

2 * * * * cd ~/cs544; git pull > /dev/null

